# server-app/deploy/oracledash.service
[Unit]
Description=__PRODUCT_NAME__ (IAM dashboard)
After=network-online.target
Wants=network-online.target

[Service]
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__CONFIG_FILE__
ExecStart=__INSTALL_DIR__/venv/bin/python __INSTALL_DIR__/app.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
